# 99-hp-wmi-permissions.rules
# /etc/udev/rules.d/99-hp-wmi-permissions.rules
#
# Grant victus-backend group write access to HP WMI fan control files.

# Find the hp-wmi device and change ownership of the pwm1_enable file
ACTION=="add", SUBSYSTEM=="hwmon", KERNELS=="hp-wmi", RUN+="/bin/sh -c 'chgrp victus-backend /sys/devices/platform/hp-wmi/hwmon/hwmon*/pwm1_enable && chmod g+w /sys/devices/platform/hp-wmi/hwmon/hwmon*/pwm1_enable'"
